# debian/postinst
#!/bin/bash
set -e
systemctl daemon-reexec
systemctl daemon-reload
systemctl enable zt-router.service
systemctl start zt-router.service
